# xe-newsserver/lib/systemd/system/fetchnews.service
[Unit]
Description=Fetch news for leafnode

[Service]
ExecStartPre=/usr/bin/nm-online -x
ExecStart=/usr/sbin/fetchnews


